# roles/weblauncher/viewer/templates/99-hdmi-sound.rules
# {{ ansible_managed }}
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/systemctl start set_hdmi_sound.service"
